# nginx.ini
[program:nginx]
command=sudo nginx 
directory=/usr/local/etc/nginx
user=root
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/nginx.log
loglevel=info
