# php-fpm.ini
[program:phpfpm]
command=sudo php-fpm
directory=/usr/var/
autorestart=true
redirect_stderr=true
stdout_logfile=/usr/var/log/php-fpm.log
loglevel=info
